# polipo.ini
[program:polipo]
command=/usr/bin/polipo -c /etc/polipo.conf

