# etc/supervisor/conf.d/aria2c.conf
[program:aria2c]
command=/usr/bin/aria2c --conf-path=/etc/aria2.conf
redirect_stderr=true
stdout_logfile=/tmp/supervisor_aria2c.log
