# etc/supervisor/conf.d/polipo.conf
[program:polipo]
command=/usr/local/bin/polipo -c /etc/polipo.conf
redirect_stderr=true
stdout_logfile=/tmp/supervisor_polipo.log
