# etc/supervisor/conf.d/ss-client.ini
[program:ss-client]
command=/usr/bin/sslocal -c /etc/ss-client.json

