# etc/supervisor/conf.d/ss-client.ini
[program:ss-client]
command=/usr/bin/sslocal -c /etc/ss-client.json
redirect_stderr=true
stdout_logfile=/tmp/ss-client.log
